# data/skeleton/etc/s6-init/run-image/service/getty/run
#!/usr/bin/execlineb -P

s6-envdir env
import -u TTY
getty -L -H DEBUG 115200 $TTY
